# roles/redisearch/templates/redisearch.service.j2
# redisearch.service
# autogenerated by Podman 2.0.6

[Unit]
Description=Podman redisearch.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/container-redisearch.pid %t/container-redisearch.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/container-redisearch.pid --cidfile %t/container-redisearch.ctr-id --cgroups=no-conmon -d --replace --rm -it --name redisearch --net host -v /var/lib/redisearch/:/data -v /etc/redisearch.conf:/etc/redis/redis.conf docker.io/redislabs/redisearch:{{ redisearch_version }} redis-server /etc/redis/redis.conf
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-redisearch.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/container-redisearch.ctr-id
PIDFile=%t/container-redisearch.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
