# roles/buildbot-worker/files/reap-dead-builds.timer
[Unit]
Description=Reap dead builds periodically

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
