# roles/buildbot-worker/files/reap-dead-builds.service
[Unit]
Description=Reap dead builds

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reap-dead-builds.py
