# roles/dehydrated/files/dehydrated.service
[Unit]
Description=Renew Let's Encrypt certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/dehydrated --cron
