# roles/linux-store-backend/templates/appstream-extractor.timer.j2
[Unit]
Description=Run appstream-extractor.service every 5 minutes

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
