# roles/linux-store-backend/templates/delete-eol-apps.timer.j2
[Unit]
Description=Run delete-eol-apps.service every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
